# app/supervisord.conf
[program:app]
command=/app-setup.sh
autorestart=false
