# conf/supervisord/mysqld.conf
[program:mysqld]
command=/bin/bash -c "exec mysqld_safe"
numprocs=1
autostart=true
autorestart=true
